# build/debian/DEBIAN/postinst
#!/bin/sh

# Change the permissions of the executables
sudo chmod +x /usr/local/sbin/rifidi/rifidi
sudo chmod 755 /usr/local/sbin/rifidi/startupscript/rifidi

# Start the init.d files
sudo ln -s /usr/local/sbin/rifidi/startupscript/rifidi /etc/init.d/rifidi
sudo update-rc.d rifidi defaults

# Start the rifidiserver
sudo /etc/init.d/rifidi start


